# systemd/moe-temp-sync.service
[Unit]
Description=Force sync MOE project to GitHub MOE-temp main
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/yuli0398/Optimization/MOE-grad-conflict-routing
Environment=MOE_TEMP_REMOTE_URL=https://github.com/lyclyq/MOE-temp.git
Environment=MOE_TEMP_BRANCH=main
EnvironmentFile=-%h/.config/moe-temp-sync.env
ExecStart=/usr/bin/env bash /home/yuli0398/Optimization/MOE-grad-conflict-routing/scripts/sync_moe_temp_force.sh
