# systemd/moe-temp-sync.timer
[Unit]
Description=Daily 11:00 force sync MOE project to GitHub MOE-temp main

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true
Unit=moe-temp-sync.service

[Install]
WantedBy=timers.target
